# services/gps/kombios-gps-service.service
[Unit]
Description=Kombi O.S. GPS Service
After=network.target

[Service]
Type=simple
User=kombios
Group=kombios
WorkingDirectory=/usr/local/bin
ExecStart=/opt/kombios/venv/bin/python /usr/local/bin/kombios-gps-service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
#Opcional: limitar recursos
#MemoryMax=200M
#CPUQuota=50%

[Install]
WantedBy=multi-user.target
